# rules/windows/process_creation/proc_creation_win_schtasks_remote_task.yml
title: Suspicious Remote Scheduled Task via schtasks
id: f3da7f14-9b54-4d6a-b85e-f3bf299c4c9b
status: experimental
description: Detects schtasks.exe creating or modifying scheduled tasks on remote systems, which may be used for lateral movement or persistence.
author: vVv-Keys
date: 2026/01/17
references:
  - https://attack.mitre.org/techniques/T1053/005/
  - https://learn.microsoft.com/windows-server/administration/windows-commands/schtasks
logsource:
  product: windows
  category: process_creation
detection:
  selection_image:
    Image|endswith:
      - '\\schtasks.exe'
      - '\\Schtasks.exe'
  selection_create:
    CommandLine|contains:
      - '/create'
      - '/change'
      - '/delete'
  selection_remote:
    CommandLine|contains:
      - '/S'
      - '\\\\'
  condition: selection_image and selection_remote and selection_create
falsepositives:
  - Legitimate remote scheduled task administration
level: medium
tags:
  - attack.execution
  - attack.lateral_movement
  - attack.t1053.005
